# servers-setup/local/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
directory=/usr/local/
command=/usr/sbin/php-fpm7.2 -F -R
autostart=true
autorestart=true
redirect_stderr=true